# _devs/duck/postgresql/etc/systemd/system/postgresql.service.d/override.conf
[Service]
Environment=PGROOT=/home/postgres
PIDFile=/home/postgres/data/postmaster.pid

ProtectHome=false
