# _devs/edgar/docker/etc/systemd/system/docker.service
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon --exec-opt native.cgroupdriver=cgroupfs -g /home/docker
